# doc/RaspberryPi3/airController.service
[Unit]
Description=AirController
After=network.target
[Service]
ExecStart=/bin/bash /home/pi/AirController/run.sh
Restart=always
RestartSec=60
User=pi
Environment=weather_api_key=<key> qingping_app_secret=<secret> sensorDataDbPassword=<password>
[Install]
WantedBy=default.target